# pig.conf
[program:pig]
command=/root/pig.py
autostart=true
autorestart=true
user=root
stderr_logfile=/var/log/pig.err.log
stdout_logfile=/var/log/pig.out.log